# .github/deployment/huducommunitybot.service
[Unit]
Description=Hudu Community Discord Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hudubot
Group=hudubot
WorkingDirectory=/opt/huducommunitybot
ExecStart=/opt/huducommunitybot/HuduCommunityBot
Restart=on-failure
RestartSec=5
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huducommunitybot

# Runtime location (installed to non-standard path via dotnet-install.sh)
Environment=DOTNET_ROOT=/usr/local/dotnet

# Bot secrets / config overrides
EnvironmentFile=/opt/huducommunitybot/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/huducommunitybot

[Install]
WantedBy=multi-user.target

